# lobordspel-server.service
[Unit]
Description=LOBordspel API Server
After=network.target mysql.service

[Service]
Type=simple
User=aideas
WorkingDirectory=/var/www/LOBordspel
ExecStart=/home/aideas/.nvm/versions/node/v24.14.0/bin/pnpm --filter server start
Restart=on-failure
RestartSec=5

Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
